# system/systemd/visaged.service.d/override.conf
[Service]
Environment=VISAGE_CAMERA_DEVICE=/dev/video-ir
Environment=VISAGE_LIVENESS_MIN_DISPLACEMENT=0.3
Environment=VISAGE_WARMUP_FRAMES=1
